# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

; ── Next.js ──────────────────────────────────────────────────────────────────
[program:nextjs]
command=node server.js
directory=/app
environment=PORT="3000",NODE_ENV="production",HOSTNAME="0.0.0.0"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Python FastAPI ───────────────────────────────────────────────────────────
[program:python]
command=uvicorn main:app --host 0.0.0.0 --port 8001 --workers 2
directory=/app/python
environment=PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
